# systemd/system/flaticron@.timer
[Unit]
Description=Run flaticron %i daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
