# systemd/system/flaticron.service
[Unit]
Description=Launch flaticron to notify of system flatpaks pending an update

[Service]
Type=oneshot
ExecStart=/usr/sbin/flaticron --system
