# systemd/user/flaticron.timer
[Unit]
Description=Run flaticron --user daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
